# systemd_units/rss-to-phlog.service
[Unit]
Description=Pull RSS and create text entries for phlog 

[Service]
Type=simple
ExecStart=/path/to/rss-to-phlog.pl